# etc/s6/haproxy/run
#!/bin/bash

haproxy -f /etc/haproxy.cfg
